# dist/systemd/system/ha-timezone-sync.service
[Unit]
Description=ha-timezone-sync

[Service]
ExecStart=/usr/bin/ha-timezone-sync
Type=oneshot

[Install]
WantedBy=multi-user.target
