# dist/systemd/user/ha-geoip.service
[Unit]
Description=ha-geoip

[Service]
ExecStart=/usr/bin/ha-geoip

[Install]
WantedBy=default.target
